# arch/systemd/system/docker-prune.service
[Unit]
Description=Prune unused Docker resources

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --force --filter "until=720h"
