# arch/systemd/system/docker-builder-prune.service
[Unit]
Description=Prune Docker builder cache

[Service]
Type=oneshot
ExecStart=/usr/bin/docker builder prune --force
