# arch/systemd/system/journal-vacuum.service
[Unit]
Description=Vacuum journal files

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=2weeks
