# arch/systemd/system/xfs-fsr.service
[Unit]
Description=Defragment XFS filesystem

[Service]
Type=oneshot
ExecStart=/usr/bin/xfs_fsr
