# arch/systemd/docker-builder-prune.service
[Unit]
Description=Weekly docker builder prune
[Service]
Type=oneshot
ExecStart=/usr/bin/docker builder prune -af
