# arch/systemd/drop-caches.service
[Unit]
Description=Drop caches
[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo 3 > /proc/sys/vm/drop_caches"
