# arch/systemd/journal-vacuum.service
[Unit]
Description=Monthly journal vacuum
[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=30d
