# arch/systemd/xfs-fsr.service
[Unit]
Description=Weekly XFS defragmentation
[Service]
Type=oneshot
ExecStart=/usr/bin/xfs_fsr
