# arch/systemd/NetworkManager.service.d/capabilities.conf
[Service]
CapabilityBoundingSet=CAP_CHOWN
