# arch/systemd-user/sway-session.target
[Unit]
Description=Sway Wayland session
Documentation=man:sway(1)
PartOf=sway.service
