# arch/systemd-user/sway-session.target
[Unit]
Description=Sway Session Target
Documentation=man:systemd.special(7)
PartOf=sway.service
