# arch/service/irqbalance.service
[Service]
ExecStart=
ExecStart=/usr/sbin/irqbalance --foreground --banmod=nvme
